# GlusterFS/assets/etc/profile.d/heketi.sh
export HEKETI_CLI_SERVER=http://k8s-1:8888
export HEKETI_CLI_USER=admin
export HEKETI_CLI_KEY=heketi